# ingar.service
[Unit]
Description=Ingar freezer alarm monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/ingar
EnvironmentFile=/home/user/.config/ingar/ingar.env
ExecStart=/usr/bin/python3 /home/user/ingar/monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
