# ingar.service
[Unit]
Description=Ingar freezer alarm monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ingar
EnvironmentFile=/home/pi/.config/ingar/ingar.env
ExecStart=/usr/bin/python3 /home/pi/ingar/monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
